# crontab
0 4 * * 1-5 curl "http://server-crawler:5000/cardapio/update" >> /outputs/out.log 2>&1
